# ansible/roles/siskin/templates/luigid.service.j2
[Unit]
Description=luigi scheduler daemon

[Service]
ExecStart=/home/{{ siskin_user }}/.local/bin/luigid

[Install]
WantedBy=multi-user.target
