# ansible/roles/siskin/templates/luigid.service.j2
[Unit]
Description=luigi scheduler daemon

[Service]
ExecStart=/home/{{ siskin_user }}/.local/bin/luigid
User={{ siskin_user }}
Group={{ siskin_group }}

[Install]
WantedBy=multi-user.target
